# opcua-simulator.service
[Unit]
Description=OPC UA Simulator Service
After=network.target

[Service]
User=kalpa
WorkingDirectory=/home/kalpa/.gemini/antigravity/scratch/opcua-simulator
ExecStart=/home/kalpa/.gemini/antigravity/scratch/opcua-simulator/venv/bin/python app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
